# deploy/systemd/samosbor-daily-review.timer
[Unit]
Description=Run Samosbor daily paper review after the MOEX futures session

[Timer]
OnCalendar=Mon..Fri *-*-* 20:45:00 UTC
Persistent=true
Unit=samosbor-daily-review.service

[Install]
WantedBy=timers.target
