# deploy/systemd/samosbor-paper-cycle.timer
[Unit]
Description=Run Samosbor paper cycle every 5 minutes during the MOEX stock main session

[Timer]
OnCalendar=Mon..Fri *-*-* 07..15:00/5:00 UTC
Persistent=true
Unit=samosbor-paper-cycle.service

[Install]
WantedBy=timers.target
