# deploy/systemd/samosbor-paper-cycle.timer
[Unit]
Description=Run Samosbor paper cycle every 5 minutes during MOEX futures session

[Timer]
OnCalendar=Mon..Fri *-*-* 06..20:00/5:00 UTC
Persistent=true
Unit=samosbor-paper-cycle.service

[Install]
WantedBy=timers.target
